# src/ca/nuci-ca.cron
MAILTO=""
1	1	*	*	*	root	/usr/share/nuci/ca/gen refresh
